# systemd/omarchy-argb.service
[Unit]
Description=Omarchy ARGB Sync Daemon
After=dbus.service framework-system.service

[Service]
Type=simple
ExecStart=/usr/local/bin/omarchy-argb daemon
Restart=on-failure

[Install]
WantedBy=default.target
